# crontab
* */6 * * * python3 /main.py >> /logfile
