# crontab
*/5 * * * * python3 /usr/src/app/example.py >> /logfile
 