# crontab
* */6 * * *  root cd /usr/src/app && python3 example.py >> /var/log/cron.log 2>&1
