# crontab
*/2 * * * *  root python3 /usr/src/app/example.py >> /var/log/cron.log 2>&1
